# config/systemd/user/xworkspace-api.service
[Unit]
Description=XWorkspace status API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/xworkspace-console/api
EnvironmentFile=-%h/.config/xworkspace/portal.env
Environment=XWORKSPACE_PORTAL_SERVICES_FILE=%h/.config/xworkspace/portal-services.json
ExecStart=/usr/bin/env go run .
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
